# etc/supervisor/conf.d/arping.conf
[program:arping]
command=bash -c 'IP=$(ip -4 addr show dev eth0 scope global | sed -nr "s/.*inet ([^/]+).*/\1/p"); /usr/bin/arping -U -c 4 $IP'
autorestart=false
autostart=true
