# etc/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf --daemonize no
autorestart=unexpected
startsecs=60
startretries=3
