# hello-world/openresty.conf
[program:openresty]
command=/opt/openresty/nginx/sbin/nginx -p /var/www/ -c /opt/openresty/nginx/config/nginx.conf  -g 'daemon off;'
autorestart=true
